# deploy/msgbot.service
[Unit]
Description=Message automation service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=%h/.local/bin/msgbot -config %h/.config/msgbot/config.toml
Restart=on-failure
RestartSec=30

[Install]
WantedBy=default.target
